# scripts/systemd/global-sentinel-dashboard.service.d/resource-limits.conf
[Service]
MemoryMax=1G
MemoryHigh=820M
CPUQuota=100%
