# scripts/systemd/global-sentinel.service
[Unit]
Description=Global Sentinel V5 Crisis Monitor
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 src/monitoring/crisis_monitor.py --repo-root /opt/global-sentinel
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=global-sentinel

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/reports
ReadWritePaths=/opt/global-sentinel/control
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
