# scripts/systemd/global-sentinel-market-warehouse-update.timer
[Unit]
Description=Run Global Sentinel Market Warehouse Daily Update

[Timer]
OnCalendar=*-*-* 06:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target

