# scripts/systemd/global-sentinel-backtest-daemon.service
[Unit]
Description=Global Sentinel Continuous Backtest Daemon
After=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/ops/continuous_backtest_daemon.py --repo-root /opt/global-sentinel
Restart=always
RestartSec=10

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/data
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target

