# scripts/systemd/global-sentinel-event-features.timer
[Unit]
Description=Run Global Sentinel Event Feature Matrix Builder

[Timer]
OnCalendar=*-*-* 06:25:00 UTC
Persistent=true

[Install]
WantedBy=timers.target

