# scripts/systemd/global-sentinel-portfolio-var.service
[Unit]
Description=Global Sentinel Portfolio VaR Calculator
After=network.target

[Service]
Type=oneshot
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/src/research/portfolio_var.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gs-portfolio-var
NoNewPrivileges=true
TimeoutStartSec=300
MemoryMax=512M

[Install]
WantedBy=multi-user.target
