# scripts/systemd/global-sentinel-portfolio-var.timer
[Unit]
Description=Run Portfolio VaR daily at 4:15 PM ET (20:15 UTC)

[Timer]
OnCalendar=*-*-* 20:15:00
Persistent=true

[Install]
WantedBy=timers.target
