# scripts/systemd/global-sentinel-canary-observability.timer
[Unit]
Description=Run Global Sentinel canary observability report periodically

[Timer]
OnBootSec=7min
OnUnitActiveSec=20min
Persistent=true

[Install]
WantedBy=timers.target
