# scripts/systemd/global-sentinel-evidence-canary.service
[Unit]
Description=Global Sentinel Evidence-Only Canary Runner
After=network.target global-sentinel.service
Wants=global-sentinel.service

[Service]
Type=oneshot
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/ops/run_evidence_only_canary.py --repo-root /opt/global-sentinel --limit 200 --window-size 25 --signal-type online_weighted_encoder

StandardOutput=journal
StandardError=journal
SyslogIdentifier=global-sentinel-evidence-canary

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/reports
ReadWritePaths=/opt/global-sentinel/telemetry
ReadWritePaths=/tmp
