# scripts/systemd/global-sentinel-canary-checkpoint.timer
[Unit]
Description=Run Global Sentinel canary stabilization checkpoint periodically

[Timer]
OnBootSec=11min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
