# scripts/systemd/global-sentinel-openclaw-research.service
[Unit]
Description=Global Sentinel OpenClaw Research Bot
After=network.target global-sentinel.service
Wants=global-sentinel.service

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
Environment=OPENCLAW_MIN_WORKERS=3
Environment=OPENCLAW_MAX_WORKERS=6
Environment=OPENCLAW_SEED_INTERVAL_SEC=60
Environment=OPENCLAW_DISPATCH_INTERVAL_SEC=30
Environment=OPENCLAW_REPLAY_EVERY_N_SEEDS=1
Environment=OPENCLAW_OUTCOME_EVERY_N_SEEDS=1
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/agent_factory.py run --bot research --config /opt/global-sentinel/config/thresholds.yaml

Restart=on-failure
RestartSec=20

StandardOutput=journal
StandardError=journal
SyslogIdentifier=global-sentinel-openclaw-research

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/reports
ReadWritePaths=/opt/global-sentinel/control
ReadWritePaths=/opt/global-sentinel/config/staging
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
