# scripts/systemd/global-sentinel-blob-health.timer
[Unit]
Description=Run Global Sentinel Blob persistence health check periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
