# scripts/systemd/gs-alpha-validator.service
[Unit]
Description=Global Sentinel Alpha Signal Validator
After=network-online.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
ExecStart=/usr/bin/python3 -u /opt/global-sentinel/src/research/alpha_validator.py
StandardOutput=journal
StandardError=journal
