# scripts/systemd/gs-compact-intents.service
[Unit]
Description=Global Sentinel: Compact Order Intents (weekly Sunday)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c '/usr/bin/python3 /opt/global-sentinel/scripts/ops/compact_intents.py >> /opt/global-sentinel/logs/compaction.log 2>&1'
User=openclaw
Group=openclaw
