# scripts/systemd/gs-daily-thesis.timer
[Unit]
Description=Daily Thesis Generator Timer

[Timer]
OnCalendar=*-*-* 12:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
