# scripts/systemd/gs-daily-performance.service
[Unit]
Description=Global Sentinel Daily Performance Attribution
After=network.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
Environment=GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/ops/daily_performance.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
