# scripts/systemd/gs-disk-cleanup.timer
[Unit]
Description=Global Sentinel Disk Cleanup Daily Timer

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
