# scripts/systemd/gs-daily-audit.service
[Unit]
Description=Global Sentinel Daily Comprehensive Audit & Auto-Fix
After=network.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
Environment=GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/ops/daily_audit.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
