# scripts/systemd/gs-earnings-calendar.timer
[Unit]
Description=Global Sentinel: Earnings Calendar Timer (daily 10:00 UTC)

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
