# scripts/systemd/gs-fundamental-valuation.timer
[Unit]
Description=Global Sentinel: Fundamental Valuation (weekdays 10:00 UTC) Timer

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
