# scripts/systemd/gs-gap-close.timer
[Unit]
Description=Save close prices at 4:05 PM ET (20:05 UTC)

[Timer]
OnCalendar=Mon..Fri *-*-* 20:05:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
