# scripts/systemd/gs-health-dashboard.timer
[Unit]
Description=Global Sentinel Health Dashboard (every 5 min)

[Timer]
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
