# scripts/systemd/gs-options-flow.service
[Unit]
Description=Global Sentinel: Options Unusual Activity Scanner (weekdays every 30min 13:30-20:00 UTC)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c "set -a && source /opt/global-sentinel/.env 2>/dev/null && set +a && /usr/bin/python3 /opt/global-sentinel/src/bridges/options_flow_scanner.py >> /opt/global-sentinel/logs/options_flow_scanner.log 2>&1"
User=openclaw
Group=openclaw
