# scripts/systemd/gs-openclaw-feed.service
[Unit]
Description=Global Sentinel: OpenClaw Bot Data Feed (weekdays every 5min market hours)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c "set -a && source /opt/global-sentinel/.env 2>/dev/null && set +a && /usr/bin/python3 /opt/global-sentinel/scripts/ops/openclaw_data_feed.py >> /opt/global-sentinel/logs/openclaw_data_feed.log 2>&1"
User=openclaw
Group=openclaw
