# scripts/systemd/gs-portfolio-optimizer.service
[Unit]
Description=Global Sentinel: Portfolio Optimization (weekdays 12:00 UTC)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c '/usr/bin/python3 -m src.research.portfolio_optimizer >> /opt/global-sentinel/logs/portfolio_optimizer.log 2>&1'
User=openclaw
Group=openclaw
