# scripts/systemd/gs-price-forecaster.service
[Unit]
Description=Global Sentinel: Price Forecasting (weekdays 12:00 & 16:00 UTC)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c '/usr/bin/python3 -m src.research.price_forecaster >> /opt/global-sentinel/logs/price_forecaster.log 2>&1'
User=openclaw
Group=openclaw
