# scripts/systemd/gs-quantum-retrain.timer
[Unit]
Description=Daily quantum retraining after market close (5 PM ET = 22:00 UTC)

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
