# scripts/systemd/gs-quantum-risk.timer
[Unit]
Description=Weekly Quantum Risk Estimation (Saturday 09:00 UTC)

[Timer]
OnCalendar=Sat *-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
