# scripts/systemd/gs-quantum-pricing.service
[Unit]
Description=Global Sentinel Quantum Option Pricer (0DTE)
After=network.target

[Service]
Type=oneshot
User=openclaw
Environment=PYTHONPATH=/opt/global-sentinel
Environment=GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/src/research/quantum_option_pricer.py
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
