# scripts/systemd/gs-strategy-orchestrator.timer
[Unit]
Description=Strategy Orchestrator Hourly
[Timer]
OnCalendar=*:00/30
Persistent=true
[Install]
WantedBy=timers.target
