# scripts/systemd/gs-sector-relative.timer
[Unit]
Description=Global Sentinel: Sector Relative Performance (weekdays every 30min 13:30-20:00 UTC) Timer

[Timer]
OnCalendar=Mon..Fri *-*-* 13..19:00,30:00 UTC
OnCalendar=Mon..Fri *-*-* 20:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
