# scripts/systemd/gs-tearsheet.service
[Unit]
Description=Global Sentinel Portfolio Tearsheet Generator
After=network-online.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
ExecStart=/usr/bin/python3 -u /opt/global-sentinel/src/research/portfolio_tearsheet.py
StandardOutput=journal
StandardError=journal
