# scripts/systemd/gs-vol-trader.service
[Unit]
Description=Global Sentinel 24/7 Volatility Opportunity Monitor
After=network-online.target global-sentinel.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 -u /opt/global-sentinel/scripts/ops/volatility_trader.py
Restart=on-failure
RestartSec=30
CPUQuota=100%
MemoryMax=512M
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gs-vol-trader

[Install]
WantedBy=multi-user.target
