# scripts/systemd/gs-vol-surface.service
[Unit]
Description=Global Sentinel: Volatility Surface Tracker (weekdays every 30min market hours)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/bin/bash -c "set -a && source /opt/global-sentinel/.env 2>/dev/null && set +a && /usr/bin/python3 /opt/global-sentinel/src/bridges/volatility_surface.py >> /opt/global-sentinel/logs/volatility_surface.log 2>&1"
User=openclaw
Group=openclaw
