# scripts/systemd/gs-watchdog.service
[Unit]
Description=Global Sentinel System Watchdog
After=network-online.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/scripts/ops/system_watchdog.py
StandardOutput=journal
StandardError=journal
