# scripts/systemd/gs-zero-dte-picker.timer
[Unit]
Description=Global Sentinel 0DTE Picker Timer (Mon-Fri 13:20 UTC / 9:20 AM ET)

[Timer]
OnCalendar=Mon..Fri *-*-* 13:20:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
