# scripts/systemd/gs-zero-dte-picker.service
[Unit]
Description=Global Sentinel 0DTE Options Strike Picker
After=network.target

[Service]
Type=oneshot
User=openclaw
WorkingDirectory=/opt/global-sentinel
Environment=PYTHONPATH=/opt/global-sentinel
Environment=GLOBAL_SENTINEL_REPO_ROOT=/opt/global-sentinel
ExecStart=/usr/bin/python3 /opt/global-sentinel/src/strategies/zero_dte_picker.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
