# scripts/systemd/healthcheck.service
[Unit]
Description=Global Sentinel Healthcheck

[Service]
Type=oneshot
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/opt/global-sentinel/.venv/bin/python scripts/healthcheck.py
