# scripts/systemd/openclaw-ops.service
[Unit]
Description=Global Sentinel V4 OpenClaw-Ops Bot (Infrastructure Orchestrator)
After=network-online.target global-sentinel.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/opt/global-sentinel/.venv/bin/python scripts/agent_factory.py list --bot ops
# In production, replace ExecStart with the OpenClaw bot runner command
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-ops

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
