# scripts/systemd/openclaw-research.service
[Unit]
Description=Global Sentinel V4 OpenClaw-Research Bot (Signal Orchestrator)
After=network-online.target global-sentinel.service
Wants=network-online.target

[Service]
Type=simple
User=gsadmin
Group=gsadmin
WorkingDirectory=/home/gsadmin/global-sentinel
EnvironmentFile=/home/gsadmin/global-sentinel/.env
ExecStart=/home/gsadmin/global-sentinel/.venv/bin/python scripts/self_improvement_loop.py
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=600
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-research

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/gsadmin/global-sentinel/logs
ReadWritePaths=/home/gsadmin/global-sentinel/reports
ReadWritePaths=/home/gsadmin/global-sentinel/config/staging
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
