# scripts/systemd/openclaw-research.service
[Unit]
Description=Global Sentinel V4 OpenClaw-Research Bot (Signal Orchestrator)
After=network-online.target global-sentinel.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/opt/global-sentinel/.venv/bin/python scripts/self_improvement_loop.py
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=600
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-research

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/reports
ReadWritePaths=/opt/global-sentinel/config/staging
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
