# scripts/systemd/self-improvement-loop.service
[Unit]
Description=Global Sentinel V4 Self-Improvement Loop (Shadow/Staging Only)
After=network-online.target openclaw-research.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/global-sentinel
EnvironmentFile=-/opt/global-sentinel/.env
Environment=PYTHONPATH=/opt/global-sentinel
ExecStart=/opt/global-sentinel/.venv/bin/python scripts/self_improvement_loop.py --repo-root /opt/global-sentinel
Restart=always
RestartSec=10
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/global-sentinel/logs
ReadWritePaths=/opt/global-sentinel/reports
ReadWritePaths=/opt/global-sentinel/control
ReadWritePaths=/opt/global-sentinel/config/staging
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
