# scripts/systemd/gs-quantum-learner.service.d/limits.conf
[Service]
MemoryMax=3G
MemoryHigh=2560M
CPUQuota=150%
Nice=10
RestartSec=120
