# scripts/systemd/gs-stop-loss.service.d/resource-limits.conf
[Service]
MemoryMax=256M
MemoryHigh=205M
CPUQuota=50%
