# scripts/systemd/global-sentinel.service.d/notifications.conf
[Service]
Environment=GS_NOTIFY_PAPER_POSITIONS=0
Environment=GS_NOTIFY_LIVE_POSITIONS=1
Environment=GS_MICROSTRUCTURE_MAX_SYMBOLS=25
Environment=GS_MICROSTRUCTURE_CACHE_TTL_SEC=3600
