# scripts/systemd/global-sentinel.service.d/resource-limits.conf
[Service]
MemoryMax=1G
MemoryHigh=920M
CPUQuota=200%
