# scripts/systemd/gs-fred-alerts.service.d/resource-limits.conf
[Service]
MemoryMax=256M
MemoryHigh=205M
CPUQuota=50%

[Unit]
StartLimitIntervalSec=600
StartLimitBurst=5
