# scripts/systemd/gs-ibkr-cpgw.service.d/resource-limits.conf
[Service]
MemoryMax=1536M
MemoryHigh=1230M
CPUQuota=100%
